# deploy/etc/systemd/system/pincabos-web.service.d/10-pincabos-web-port.conf
[Service]
Environment=PINCABOS_WEB_PORT=5055
