# deploy/etc/systemd/system/pincabos-web.service.d/10-pincabos-legacy.conf
[Unit]
Documentation=https://pincabos.cc/
