# deploy/etc/systemd/system/getty@tty1.service.d/10-pincabos-final-graphical-guard.conf
[Unit]
After=pincabos-final-graphical-guard.service
Wants=pincabos-final-graphical-guard.service

[Service]
ExecStartPre=-/usr/local/sbin/pincabos-final-graphical-guard.sh --from-getty
