# deploy/etc/systemd/system/pincabos-switch-graphical-vt.service
[Unit]
Description=PinCabOS Switch to Graphical VT
Documentation=https://pincabos.cc/
After=lightdm.service display-manager.service graphical.target pincabos-final-graphical-guard.service
Wants=lightdm.service

[Service]
Type=oneshot
ExecStart=/bin/bash -lc 'for i in $(seq 1 60); do [ -S /tmp/.X11-unix/X0 ] && break; sleep 1; done; command -v chvt >/dev/null 2>&1 && chvt 7 || true; pco_go "pincabos-vpinfe restart deferred until final reboot"'
RemainAfterExit=yes
TimeoutStartSec=90

[Install]
WantedBy=multi-user.target graphical.target
