# deploy/etc/systemd/system/pincabos-web.service
[Unit]
Description=PinCabOS WebApp Legacy Alias
After=pincabos-webapp.service
Requires=pincabos-webapp.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/true

[Install]
WantedBy=multi-user.target
