# deploy/etc/systemd/system/pincabos-webapp.service.d/20-sudo-privilege-transition.conf
# PinCabOS WebApp:
# WebApp runs as pinball on localhost:5055.
# Root operations remain limited by /etc/sudoers.d/pincabos-web.
# Do not restrict CapabilityBoundingSet: whitelisted sudo commands such as
# update-gpu-drivers.sh need normal root capabilities for APT and state flags.
[Service]
NoNewPrivileges=no
RestrictSUIDSGID=no
