# deploy/etc/systemd/system/pincabos-webapp.service.d/40-waitress.conf
[Service]
ExecStart=
ExecStart=/opt/pincabos/web/.venv/bin/waitress-serve --host=0.0.0.0 --port=80 --threads=8 --channel-timeout=120 app:app
